# systemd/mariadbBackup.service
[Unit]
Description="Backup mariadb database"
Requires=docker.service
After=docker.service

[Service]
ExecStart=/root/.export_mariadb.sh
