# systemd/mariadbBackup.timer
[Unit]
Description="Backup mariadb database"

[Timer]
Unit=mariadbBackup.service
OnCalendar=*-*-* 02:00:00

[Install]
WantedBy=timers.target